# worker/crontab
*/5 * * * * root python fetch-data.sh >> /code/fetch-data.log 2>&1

